# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[inet_http_server]
port = *:9001

[supervisord]
nodaemon=true
logfile=%(here)s/logs/supervisord.log
logfile_maxbytes=5MB
logfile_backups=5
loglevel=info
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ── API Gunicorn + Gevent ─────────────────────────────────────────
[program:app]
command=gunicorn app.main:app -k uvicorn.workers.UvicornWorker -b 0.0.0.0:8000
directory=%(here)s
autorestart=true
startretries=5
startsecs=5
stdout_logfile=%(here)s/logs/app.log
stderr_logfile=%(here)s/logs/app.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=info
environment=PYTHONPATH="%(here)s"

; ── Scraper Manager (tous pays) ───────────────────────────────────
; Ajouter un pays = créer un fichier dans app/worker/sources/
; Aucune modification de ce fichier requise.
[program:scraper_manager]
command=python scraper_manager.py
directory=%(here)s/app/worker
autorestart=true
startretries=10
startsecs=3
stdout_logfile=%(here)s/logs/scraper_manager.log
stderr_logfile=%(here)s/logs/scraper_manager.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=100
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=info
environment=PYTHONPATH="%(here)s/../.."
